# host-config/laptop/com.magnusgille.brew-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.magnusgille.brew-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/magnus/.local/bin/grimnir-brew-update.sh</string>
    </array>

    <!-- Weekly: Sunday 11:00 local. launchd runs a missed calendar job on next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/magnus/.local/share/grimnir-maintenance/logs/brew-update.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/magnus/.local/share/grimnir-maintenance/logs/brew-update.log</string>
</dict>
</plist>
