# systemd/grimnir-validate.service
[Unit]
Description=Grimnir registry validation — compare services.json against live state
After=network-online.target munin-memory.service
Wants=network-online.target

[Service]
Type=oneshot
User=grimnir
WorkingDirectory=/srv/grimnir/control-plane
LoadCredential=munin-api-key:/etc/grimnir/credentials/munin-api-key
ExecStart=/usr/bin/bash scripts/generate-architecture.sh --munin-token-file %d/munin-api-key --validate

# Timeout — cross-host SSH checks can be slow
TimeoutStartSec=120

# Sandboxing — read-only, needs network for SSH and Munin HTTP
ProtectSystem=strict
ReadOnlyPaths=/srv/grimnir
# Exception (#33): let the validate run self-heal ONLY the deploy marker so
# Heimdall's drift detector doesn't false-flag grimnir after a session pulls the
# tree forward outside a deploy. The '-' prefix makes it optional, so a checkout
# without a marker yet still starts cleanly. Overrides ReadOnlyPaths/ProtectHome
# for this single file only.
ReadWritePaths=-/srv/grimnir/control-plane/.deployed-commit
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
PrivateDevices=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
ProtectClock=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
SystemCallArchitectures=native

[Install]
WantedBy=default.target
