# systemd/grimnir-maintenance-deps.timer
[Unit]
Description=Grimnir weekly npm dependency report

[Timer]
# Monday early morning, off-peak and clear of the Sunday 03:00 security scan.
OnCalendar=Mon *-*-* 02:10:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
